# pkg/centos/after-install.sh
/sbin/chkconfig --add logstash

chown -R logstash:logstash /opt/logstash
chown logstash /var/log/logstash
chown logstash:logstash /var/lib/logstash
